# root/etc/services.d/web/run
#!/usr/bin/execlineb -P

with-contenv

backtick -D secret PORTUS_SECRET_KEY_BASE
  { cd /portus bundle exec rake secret }
  cd /portus bundle exec puma -C /usr/local/etc/portus/puma.rb
