# root/etc/services.d/crono/run
#!/usr/bin/execlineb -P

with-contenv

cd /portus bundle exec crono -L STDOUT -e production
