# root/etc/services.d/crono/run
#!/usr/bin/with-contenv sh

set +e

cd /portus

exec /usr/local/bundle/bin/bundle exec crono
